# basic/config/cronPS/all.cron
10 * * * * sh ../cronPSConnector/all.sh "Push"
*/2 * * * * sh ../cronPSConnector/all.sh "Realtime"