# basic/config/cronPS/private.cron
5 * * * * sh ../cronPSConnector/private.sh "Push"
*/1 * * * * sh ../cronPSConnector/private.sh "Realtime"